# .docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=10

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=20

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --tries=3
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/supervisor/worker.log
stderr_logfile=/var/log/supervisor/worker_error.log
priority=30

[program:vite]
command=npm run dev
directory=/var/www
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/vite.log
stderr_logfile=/var/log/supervisor/vite_err.log
priority=40